# wms-app.service
[Unit]
Description=WMS Next.js Application
After=network.target

[Service]
Type=simple
User=wms
Group=wms
WorkingDirectory=/home/wms/wms-app/.next/standalone
Environment="NODE_ENV=production"
Environment="BASE_PATH=/WMS"
Environment="PORT=3000"
Environment="HOSTNAME=0.0.0.0"
ExecStart=/usr/bin/node /home/wms/wms-app/.next/standalone/server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=wms-app

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/wms/wms-app

[Install]
WantedBy=multi-user.target